# go-scoreboard/init/scoreboard.init
#!/sbin/openrc-run

supervisor=supervise-daemon